# files/cozy-controller.conf
[program:cozy-controller]
command=/usr/local/lib/node_modules/cozy-controller/bin/cozy-controller -u --auth --per 755
environment=NODE_ENV="production"
user=root
