# files/cozy-indexer.conf
[program:cozy-indexer]
user = cozy
command = /usr/local/cozy-indexer/cozy-data-indexer/virtualenv/bin/python /usr/local/cozy-indexer/cozy-data-indexer/server.py
